# .config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q ~/.ssh/github < /dev/null
Name=ssh-add
Type=Application
